# roles/nginx/templates/certbot-update.timer.j2
[Unit]
Description=Certbot update timer
Requires=certbot-update.service

[Timer]
Unit=certbot-update.service
OnCalendar=monthly

[Install]
WantedBy=timers.target
