# roles/cf-ddns/templates/cf-ddns.timer.j2
[Unit]
Description=CloudFlare DDNS timer
Requires=cf-ddns.service

[Timer]
Unit=cf-ddns.service
OnCalendar=daily

[Install]
WantedBy=timers.target
